# etc/rc.d/rc.local
#!/bin/sh

ntpdate -u pool.ntp.org
/root/term1.sh &
